# cfg/supervisord_prod.cfg
[program:act2disp]
command=/usr/bin/python /opt/monpanel.com/hub/bin/act2disp.py
directory=/opt/monpanel.com/hub/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/log/act2disp.err.log
stdout_logfile=/opt/monpanel.com/hub/log/act2disp.out.log
user=root
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub'

[program:pulldata]
command=/usr/bin/python /opt/monpanel.com/hub/bin/pulldata.py
directory=/opt/monpanel.com/hub/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/log/pulldata.err.log
stdout_logfile=/opt/monpanel.com/hub/log/pulldata.out.log
user=root
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub'

[program:pushdata]
command=/usr/bin/python /opt/monpanel.com/hub/bin/pushdata.py
directory=/opt/monpanel.com/hub/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/log/pushdata.err.log
stdout_logfile=/opt/monpanel.com/hub/log/pushdata.out.log
user=mp
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub'

[program:rttmon]
command=/usr/bin/python3 /opt/monpanel.com/hub/bin/rttmon.py
directory=/opt/monpanel.com/hub/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/log/rttmon.err.log
stdout_logfile=/opt/monpanel.com/hub/log/rttmon.out.log
user=mp
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub'


