# cfg/supervisord_prod.cfg
[program:act2disp]
command=/usr/bin/python /opt/monpanel.com/hub/prod/bin/act2disp.py
directory=/opt/monpanel.com/hub/prod/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/prod/log/act2disp.err.log
stdout_logfile=/opt/monpanel.com/hub/prod/log/act2disp.out.log
user=root
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub/prod'

[program:pulldata]
command=/usr/bin/python /opt/monpanel.com/hub/prod/bin/pulldata.py
directory=/opt/monpanel.com/hub/prod/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/prod/log/pulldata.err.log
stdout_logfile=/opt/monpanel.com/hub/prod/log/pulldata.out.log
user=root
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub/prod'

[program:pushdata]
command=/usr/bin/python /opt/monpanel.com/hub/prod/bin/pushdata.py
directory=/opt/monpanel.com/hub/prod/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/prod/log/pushdata.err.log
stdout_logfile=/opt/monpanel.com/hub/prod/log/pushdata.out.log
user=mp
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub/prod'

[program:rttmon]
command=/usr/bin/python /opt/monpanel.com/hub/prod/bin/rttmon.py
directory=/opt/monpanel.com/hub/prod/tmp
autostart=true
autorestart=true
startretries=1
numprocs=1
stderr_logfile=/opt/monpanel.com/hub/prod/log/rttmon.err.log
stdout_logfile=/opt/monpanel.com/hub/prod/log/rttmon.out.log
user=mp
environment=PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin':HUB_HOME='/opt/monpanel.com/hub/prod'

[group:hubsrv]
programs=pulldata,pushdata,rttmon,act2disp

